# files/etc/uci-defaults/999-QOL_config
uci -q batch << EOI

# Wireless default settings
set wireless.radio0.country='US'
set wireless.radio1.country='US'
set wireless.radio2.country='US'
set wireless.radio0.disabled=1
set wireless.radio1.disabled=1
set wireless.radio2.disabled=0

# Network settings
set network.globals.packet_steering='0'
set network.@device[0].flow_offloading='0'
set network.@device[0].flow_offloading_hw='0'

EOI
